# systemd/system/ttrss-bak.service
[Unit]
Description=Backup TT-RSS user data

[Service]
ExecStart=/usr/local/bin/ttrss-bak
Type=oneshot
TimeoutSec=30
