# systemd/user/arbie.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=*-*-* 08:00:00