# systemd/user/sleep-notify.timer
[Unit]
Description=Sleep time notification

[Timer]
OnCalendar=*-*-* 18:00:00
Unit=sleep-notify.service

[Install]
WantedBy=default.target
