# systemd/user/yt-dlp.timer
[Unit]
Description=Youtube Downloader 

[Timer]
OnCalendar=hourly
Persistent=true
Unit=yt-dlp.service

[Install]
WantedBy=default.target
