# systemd/user/yt-dlp.service
[Unit]
Description=Youtube Downloader
Wants=default.target

[Service]
Type=simple
ExecStart=yt-dlp
WorkingDirectory=%h/Videos
